# supervisord.conf
[supervisord]
nodaemon=true

[program:motorsim]
command=/bin/bash -c "cd /epics/iocs/motorsim && ./st.cmd"

[program:adsim]
command=/bin/bash -c "cd /epics/iocs/adsim && ./st.cmd"
